# internal/idle/scripts/mutapod-idle-check.service
[Unit]
Description=Mutapod idle shutdown check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/mutapod-idle-check
